# files/supervisor/filebeat.ini
[program:filebeat]
user = root
autostart = true
command = /opt/filebeat/filebeat -e -c /opt/filebeat/filebeat.yml
stdout_logfile = /var/log/monitor/filebeat-common.log
stderr_logfile = /var/log/monitor/filebeat-error.log
