# files/supervisor/metricbeat.ini
[program:metricbeat]
user = root
autostart = true
command = /opt/monitor/metricbeat/metricbeat -e -c /opt/monitor/metricbeat/metricbeat.yml
stdout_logfile = /var/log/monitor/metricbeat/common.log
stderr_logfile = /var/log/monitor/metricbeat/error.log
