# files/supervisor/packetbeat.ini
[program:packetbeat]
user = root
autostart = true
command = /opt/monitor/packetbeat/packetbeat -e -c /opt/monitor/packetbeat/packetbeat.yml
stdout_logfile = /var/log/monitor/packetbeat/common.log
stderr_logfile = /var/log/monitor/packetbeat/error.log
